# wn-ac1600dgr2/GPL_RELEASE/ISD2/TARGET_RAMDISK/MIPS32_STAGE/etc/inittab
console::sysinit:/sbin/config_init
ttyS0::respawn:/sbin/config_term
::ctrlaltdel:/etc/reboot.sh
::shutdown:/etc/reboot.sh
::shutdown:/etc/reboot.sh /sbin/swapoff -a
