# wn-ac1600dgr2/GPL_RELEASE/ISD2/TARGET_RAMDISK/MIPS32_STAGE/etc/rc.d/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here.


